# lfbfi.service
[Unit]
After=dhcpcd.service

[Service]
ExecStart=/bin/fb

[Install]
WantedBy=default.target
